# infra/systemd/budgetbook-backup.timer
[Unit]
Description=Run BudgetBook SQLite backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10m
Unit=budgetbook-backup.service

[Install]
WantedBy=timers.target
